# build/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:backend]
command=/app/cms
environment=ADDR=0.0.0.0:8080

[program:frontend]
command=npm run start
environment=NODE_ENV=production

[program:reverse-proxy]
command=caddy run --config /etc/caddy/Caddyfile
environment=
