# build/supervisord.conf
[supervisord]
nodaemon=true
user=cms

[program:backend]
command=cms
environment=ADDR=0.0.0.0:8080

[program:frontend]
command=node server.js
environment=

[program:reverse-proxy]
command=caddy run --config /etc/caddy/Caddyfile
environment=
